# .config/systemd/user/mediasync.service
[Unit]
Description=sync media
 
[Service]
Type=oneshot
ExecStart=/bin/bash /home/will/scripts/syncer.sh
